# assets/templates/jacker-compose-reload.timer.template
[Unit]
Description=Refresh images and update Jacker containers
Requires=jacker-compose.service
After=jacker-compose.service

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target