# assets/templates/jacker-compose-reload.service.template
[Unit]
Description=Refresh images and update Jacker containers

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload-or-restart jacker-compose.service
